# deploy/systemd/market-brief.service
[Unit]
Description=Generate due Market Brief reports
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=marketbrief
Group=marketbrief
WorkingDirectory=/opt/market-brief
Environment=HOME=/home/marketbrief
Environment=TZ=America/Chicago
Environment=MARKET_BRIEF_PROJECT=/opt/market-brief
Environment=MARKET_BRIEF_RUNNER=/opt/market-brief/run_portable.sh
EnvironmentFile=-/etc/market-brief.env
ExecStart=/opt/market-brief/schedule.sh
TimeoutStartSec=45min
Nice=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/market-brief/data /opt/market-brief/out /opt/market-brief/logs /opt/market-brief/state /home/marketbrief

[Install]
WantedBy=multi-user.target
